# deploy/systemd/bookedai-boot-reconcile.service
[Unit]
Description=Reconcile BookedAI stack after boot
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/dovanlong/BookedAI
ExecStart=/bin/bash /home/dovanlong/BookedAI/scripts/reconcile_boot_state.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
