# deploy/systemd/bookedai-cloudflare-dns.service
[Unit]
Description=Update BookedAI Cloudflare DNS records on boot
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/dovanlong/BookedAI
ExecStart=/bin/bash /home/dovanlong/BookedAI/scripts/update_cloudflare_dns_records.sh

[Install]
WantedBy=multi-user.target
